# debian/inittab
::sysinit:/etc/init.d/rcS
::respawn:/bin/sh
::shutdown:/etc/init.d/rc 0
